# etc/daemon.ini
[program:check_sender_access]
command=/opt/zmbr_check_sender/check_sender_access.py
process_name=zimbra_check_sender
priority=1
redirect_stderr=true
stdout_logfile=/var/log/check_sender_access.log
